# cfg/systemd/system/machines.target
[Unit]
Description=Containers
Documentation=man:systemd.special(7)
Before=multi-user.target

[Install]
WantedBy=multi-user.target
