# cfg/systemd/system/remote-fs-pre.target
[Unit]
Description=Remote File Systems (Pre)
Documentation=man:systemd.special(7)
RefuseManualStart=yes
