# cfg/systemd/system/sigpwr.target
[Unit]
Description=Power Failure
Documentation=man:systemd.special(7)
