# cfg/systemd/system/sleep.target
[Unit]
Description=Sleep
Documentation=man:systemd.special(7)
DefaultDependencies=no
RefuseManualStart=yes
StopWhenUnneeded=yes
