# cfg/systemd/system/swap.target
[Unit]
Description=Swap
Documentation=man:systemd.special(7)
