# cfg/systemd/system/timers.target
[Unit]
Description=Timers
Documentation=man:systemd.special(7)

DefaultDependencies=no
Conflicts=shutdown.target
