# cfg/systemd/system/umount.target
[Unit]
Description=Unmount All Filesystems
Documentation=man:systemd.special(7)
DefaultDependencies=no
RefuseManualStart=yes
